# another-packages/shaper-2.2.12-0.7.3/debian/shaper.cron.d
# This is needed in order for the TIME parameter in /etc/shaper/cbq-* to
# work.
*/10 * * * *	root	[ -d /etc/shaper ] && /etc/init.d/shaper timecheck
