# another-packages/sarg-2.3.1/debian/sarg-reports.cron.daily
#!/bin/sh

if [ -x /usr/sbin/sarg-reports ]; then
  /usr/sbin/sarg-reports daily
fi
